# aws-cli/put_httpd_metric_data.conf
start on runlevel [23]
stop on shutdown
# Respawn in case of a crash, with default parameters
respawn

script
exec /path/to/put_httpd_metric_data.sh
end script

